# scripts/sweety.service
[Unit]
Description=Sweety High-Performance Web Server
Documentation=https://github.com/your-org/sweety
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=www-data
Group=www-data

# 可执行文件路径（安装后修改）
ExecStart=/usr/local/bin/sweety run --config /etc/sweety/sweety.toml
ExecReload=/usr/local/bin/sweety reload --config /etc/sweety/sweety.toml
ExecStop=/bin/kill -SIGTERM $MAINPID

# 工作目录
WorkingDirectory=/etc/sweety

# 重启策略：非正常退出时自动重启，5 秒后重试
Restart=on-failure
RestartSec=5s

# 文件描述符上限（对标 Nginx worker_rlimit_nofile）
LimitNOFILE=1048576
LimitNPROC=65536

# 安全加固（可选，根据实际需求调整）
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=yes

# 日志输出到 systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sweety

# 超时设置
TimeoutStartSec=30
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
